# chef/cookbooks/metasploitable/files/sinatra/sinatra.conf
description 'Run vulnerable Sinatra'
author 'metasploitable3'

start on runlevel [2345]
exec "/opt/sinatra/server"
