# chef/cookbooks/metasploitable/files/sinatra/sinatra.conf
description 'Run vulnerable Sinatra'
author 'metasploitable3'

start on runlevel [2345]
chdir /opt/sinatra/
script
  echo "Starting sinatra"
  /opt/sinatra/server
end script