# chef/cookbooks/docker/templates/default/systemd/docker.socket.erb
[Unit]
Description=Docker Socket for the API
PartOf=<%= @docker_name %>.service

[Socket]
ListenStream=<%= @docker_socket %>
SocketMode=0660
SocketUser=root
SocketGroup=docker

[Install]
WantedBy=sockets.target
