# chef/cookbooks/metasploitable/files/readme_app/readme_app.conf
description 'Run ReadMe App'
author 'metasploitable3'

start on runlevel [2345]

script
    sudo -u chewbacca /opt/readme_app/start.sh
end script
