# chef/cookbooks/metasploitable/files/proftpd/hosts_renewer.conf
description 'Run hosts renewer'
author 'metasploitable3'

start on runlevel [2345]
script
  ruby /opt/proftpd/hosts_renewer.rb
end script