# chef/cookbooks/mysql/templates/default/systemd/mysqld.service.erb
[Unit]
Description=mysql_service[mysql-<%= @config.instance %>]
After=syslog.target
After=network.target

[Service]
Type=simple
User=<%= @config.run_user %>
Group=<%= @config.run_group %>
ExecStart=<%= @mysqld_bin %> --defaults-file=<%= @etc_dir %>/my.cnf --basedir=<%= @base_dir %>
ExecStartPost=/usr/libexec/mysql-<%= @config.instance %>-wait-ready $MAINPID
TimeoutSec=300
PrivateTmp=true

[Install]
WantedBy=multi-user.target
